# scripts/systemd/xcash-rpc-wallet-turbo-tx.timer
[Unit]
Description=xcash-rpc-wallet-turbo-tx timer

[Timer]
OnBootSec=20sec

[Install]
WantedBy=timers.target

