# scripts/systemd/xcash-rpc-wallet.timer
[Unit]
Description=xcash-rpc-wallet timer

[Timer]
OnBootSec=20sec

[Install]
WantedBy=timers.target

