# scripts/systemd/xcash-dpops.timer
[Unit]
Description=xcash-dpops timer

[Timer]
OnBootSec=30sec

[Install]
WantedBy=timers.target

